# .deploy-stage/tvserverd.service
[Unit]
Description=TV Kit API and RÚV scheduler
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=olafurbui
Group=olafurbui
WorkingDirectory=/opt/tv-kit
EnvironmentFile=/etc/tv-kit/tv-kit.env
Environment=HOME=/var/lib/tv-kit
ExecStart=/opt/tv-kit/bin/tvserverd
Restart=always
RestartSec=2
StateDirectory=tv-kit
LogsDirectory=tv-kit
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/tv-kit /var/log/tv-kit
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tvserverd

[Install]
WantedBy=multi-user.target
