# ops/systemd/user/tv-dashboard-dev.service
[Unit]
Description=TV Kit dashboard Vite development server
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/tv-kit-dev
ExecStart=%h/.bun/bin/bun run dev:dashboard
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
