# ops/systemd/user/tv-kiosk.service
[Unit]
Description=TV Kit Chrome kiosk
After=graphical-session.target tv-dashboard.service
Wants=tv-dashboard.service

[Service]
Type=simple
Environment=DISPLAY=:1
Environment=WAYLAND_DISPLAY=wayland-0
# Kill lingering Chrome processes before every (re)start so old instances
# never pile up or steal the kiosk window.
ExecStartPre=-/usr/bin/flatpak kill com.google.Chrome
ExecStartPre=-/usr/bin/pkill -9 -f com.google.Chrome
ExecStart=/usr/bin/flatpak run com.google.Chrome --kiosk --autoplay-policy=no-user-gesture-required --noerrdialogs --disable-session-crashed-bubble --disable-infobars --start-fullscreen http://127.0.0.1:3111
ExecStopPost=-/usr/bin/pkill -9 -f com.google.Chrome
Restart=always
RestartSec=3
TimeoutStopSec=10

[Install]
WantedBy=default.target
