# ops/systemd/user/tv-remote-dev.service
[Unit]
Description=TV Kit remote Vite development server
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/tv-kit-dev
ExecStart=%h/.bun/bin/bun run dev:remote
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
