# ops/systemd/user/tvserverd-radioscraper.service
[Unit]
Description=Refresh and validate Homeboard radio stations
Wants=network-online.target
After=network-online.target tvserverd.service

[Service]
Type=oneshot
TimeoutStartSec=2min
WorkingDirectory=%h/.tv-kit/src
EnvironmentFile=%h/.tv-kit/src/.env
EnvironmentFile=-%h/.tv-kit/env
Environment=PATH=%h/.bun/bin:/usr/local/bin:/usr/bin
ExecStart=/usr/bin/env bun apps/server/src/radioscraper.ts
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
SyslogIdentifier=tvserverd-radioscraper
