# ops/systemd/user/tvserverd.service
[Unit]
Description=Homeboard TV server daemon
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Projects/tv-kit
EnvironmentFile=%h/Projects/tv-kit/.env
Environment=PATH=%h/.bun/bin:/usr/local/bin:/usr/bin
ExecStart=/usr/bin/env bun apps/server/src/index.ts
Restart=always
RestartSec=2
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
SyslogIdentifier=tvserverd

[Install]
WantedBy=default.target
